# crontab
*/1 * * * * /usr/bin/php -f /home/www/public_html/index.php